# src/install_systemd/beereader.timer
[Unit]
Description=Check for available BT devices and attempt connection.

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target

