# src/install_systemd/beereader.timer
[Unit]
Description=Run database value collection and send data mail.

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target

