# src/install_systemd/databasecsv.timer
[Unit]
Description=Run database value collection and send data mail.

[Timer]
OnBootSec=1min
OnUnitActiveSec=24h

[Install]
WantedBy=timers.target

